# docker/php-worker/supervisord.d/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan rabbitmq:consume
autostart=true
autorestart=true
numprocs=1
user=laradock
